# deploy/crontab
MAILTO=""
10 * * * * /opt/util/check.py > /proc/1/fd/1 2>/proc/1/fd/2